# theme-switcher/.config/systemd/user/theme-github-contributions.timer
[Unit]
Description=Hourly theme dashboard GitHub contribution refresh

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
Unit=theme-github-contributions.service

[Install]
WantedBy=timers.target
